# tools/com.example.checkapps.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.checkapps</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>sleep 300 &amp;&amp; $HOME/Technology/tools/check-apps.sh notify</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
